# root/etc/services.d/queue/run
#!/usr/bin/with-contenv bash
exec s6-setuidgid abc /usr/bin/php /app/src/artisan queue:listen --sleep=3 --tries=3 --queue="default"